# roles/workstation/files/dotfiles/autostart/protonvpn.desktop
[Desktop Entry]
Exec=/usr/bin/protonvpn-cli c -f
