# roles/workstation/files/dotfiles/autostart/protonvpn.desktop
[Desktop Entry]
Categories=Utility;
StartupNotify=false
Terminal=true
Type=Application
Exec=/usr/bin/protonvpn-cli c -f
Name=ProtonVPN Autoconnect
